# acme.crontab
# minute	hour	mday	month	wday	command
30	9	*	*	*	/usr/local/bin/acmecron.sh
