# certmgr.crontab
# minute	hour	mday	month	wday	command
30	10	*	*	*	/usr/local/bin/updatecerts.sh
